# dist/mastotg.service
[Unit]
Description=mastotg round

[Service]
Type=oneshot
# Replace the option values with yours
ExecStart=/usr/local/bin/mastotg -i query-fetch -o tg-send -s social.myl.moe -u myl --tg-chan myl7s -f /srv/mastotg/mastotg.json
# Put the envs here. It should not be so bad if the unit file is in mode 400.
# For more safety, refer to LoadCredential.
Environment=TELOXIDE_TOKEN=YOUR_TG_BOT_TOKEN
# Create a standalone user for the program to gain more safety
#User=
#Group=
